# target/linux/brcm2708/base-files/etc/crontabs/root
@reboot /usr/sbin/logrotate /etc/logrotate.conf
@hourly /usr/sbin/logrotate /etc/logrotate.conf

